# Docker_dind/docker
#!/sbin/openrc-run

start() {
  sleep 2
  /usr/local/bin/dockerd &
}

stop () {
  kill $(pidof dockerd)
}
